# rustic/.config/systemd/user/backup.timer
[Unit]
Description="Daily Backup"

[Timer]
Unit=backup.service
OnCalendar=10:00
Persistent=true

[Install]
WantedBy=timers.target

