# rustic/.config/systemd/user/backup.service
[Unit]
Description=Run Backup Script

[Service]
Type=oneshot
ExecStart=/home/giuseppe/.local/scripts/backup.sh
